# deploy/certificados/gunicorn-certificados.service
# /etc/systemd/system/gunicorn-certificados.service
# Ajustar User= y rutas. Texto largo de actividad en ACTIVIDAD_TXT.

[Unit]
Description=Gunicorn certificados.relaticpanama.org
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/certificados
Environment="PATH=/var/www/certificados/venv/bin"
Environment="ACTIVIDAD_TXT=Gestión Editorial en Revistas Científicas, convenio vigente No. 100-001 del 20 de octubre de 2025 entre UMLA y Relatic Panamá."
ExecStart=/var/www/certificados/venv/bin/gunicorn -w 2 -b 127.0.0.1:8000 app.app:app

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
